# dest/systemd/unclutter.service
[Unit]
Description=Starts unclutter

[Service]
ExecStart=unclutter -idle 2 -root
KillMode=process
