# dest/systemd/xbindkeys.service
[Unit]
Description=Starts xbindkeys

[Service]
ExecStart=xbindkeys --nodaemon
KillMode=process
