# systemd/steam.service
[Unit]
Description=Starts Steam headless

[Service]
ExecStart=steam -noasync -no-browser -silent
Restart=always
RestartSec=5s
StandardOutput=null
