# systemd/weatherdump.timer
[Unit]
Description=Hourly local weather dump

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=default.target
