# systemd/xpra-client.service
[Unit]
Description=Connects to xpra display

[Service]
ExecStart=xpra attach :10 --dpi=96
Restart=always
RestartSec=5s
