# systemd/xsuspender.service
[Unit]
Description=Starts xsuspender

[Service]
ExecStart=xsuspender
KillMode=process
